# toutiao-backend/scripts/toutiao_app_supervisor.conf
[program:toutiao-app]
command=/home/python/toutiao-backend/scripts/toutiao_app.sh
directory=/home/python/toutiao-backend
user=python
autorestart=true
redirect_stderr=false
loglevel=info
stopsignal=KILL
stopasgroup=true
killasgroup=true